# root/etc/services.d/rsyslog/run
#!/usr/bin/with-contenv bash

if test "${DCO_RSYSLOG_START}" = "1"; then
    if test "${DCO_RSYSLOG_REMOTE_HOST}" != "null"; then
        echo "# AUTOGENERATED FROM remote.template => DO NOT MODIFY HERE" >/etc/rsyslog.d/remote.conf
        cat /etc/rsyslog.d/remote.template |sed "s/~~~DCO_RSYSLOG_REMOTE_HOST~~~/${DCO_RSYSLOG_REMOTE_HOST}/g" |sed "s/~~~DCO_RSYSLOG_REMOTE_PORT~~~/${DCO_RSYSLOG_REMOTE_PORT}/g" >>/etc/rsyslog.d/remote.conf
    else
        rm -f /etc/rsyslog.d/remote.conf
    fi
    exec /sbin/rsyslogd -n
else
    exec /bin/sleep infinity
fi
